# deploy/supervisor/ai-service.conf
[program:shiforeyting-ai]
directory=/var/www/shiforeyting/ai-service
command=/var/www/shiforeyting/ai-service/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001 --workers 2
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/shiforeyting/storage/logs/ai-service.log

